# spark-standalone/sparkslave.conf
[program:spark-slave]
command=/bin/bash -c "/opt/spark-1.2.0-bin-hadoop1/sbin/start-slave.sh 1 spark://${HOSTNAME}:7077"
