# spark-standalone/sparkslave.conf
[program:spark-slave]
command=/home/spark/spark/sbin/start-slave.sh 1 spark://$hostname:7077
