# spark-standalone/sparkslave.conf
[program:spark-slave]
command=/bin/bash -c "/home/spark/spark/sbin/start-slave.sh 1 spark://${HOSTNAME}:7077"
