# spark-standalone/sparkmaster.conf
[program:spark-master]
command=/opt/spark-1.2.0-bin-hadoop1//sbin/start-master.sh
