# spark-standalone/sparkmaster.conf
[program:spark-master]
command=/home/spark/spark/sbin/start-master.sh
